# net-wireless/stlc45xx-cal/files/stlc45xx.rules
# do not edit this file, it will be overwritten on update

# /etc/udev/rules/40-stlc45xx.rules:  calibrating stlc45xx

SUBSYSTEM=="net", DRIVERS=="cx3110x", ACTION=="add", RUN+="stlc45xx-cal"
