# net-wireless/tablet-wireless/files/tablet-wireless.init
#!/sbin/runscript

depend() {
	before net
}

start() {
	ebegin "Loading cx3110x module"
	insmod /mnt/initfs/lib/modules/`uname -r`/umac.ko
	modprobe cx3110x || modprobe cx3110x_mt
	eend $?
	ebegin "Calibrating wlan.."
	chroot /mnt/initfs /usr/bin/wlan-cal
	eend $?
}

stop() {
	ebegin "Unloading cx3110x module"
	rmmod cx3110x || rmmod cx3110x_mt
	eend $?
}
