# net-wireless/tablet-wireless/files/tablet-wireless.init
#!/sbin/runscript

depend() {
	before net
}

start() {
	einfo "Loading cx3110x module"
	insmod /mnt/initfs/lib/modules/`uname -r`/umac.ko
	modprobe cx3110x
	einfo "Calibrating wlan.."
	chroot /mnt/initfs /usr/bin/wlan-cal
}
