# net-wireless/cx3110x/files/cx3110x.rules
# do not edit this file, it will be overwritten on update

# /etc/udev/rules/40-cx3110x.rules:  calibrating cx3110x

DRIVER=="wlan-omap", ACTION=="add", RUN+="cx3110x-cal"
